# laravel-api/docker/fusionpbx-docker/supervisord.conf
[supervisord]
nodaemon=true

[program:php7-fpm]
command=/usr/sbin/php-fpm7.3 --nodaemonize -c /etc/php/7.3/fpm/php-fpm.conf

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

[program:freeswitch]
autorestart=true
command=/bin/sh /usr/bin/start-freeswitch.sh